# systemd/calibre-server.service
[Unit]
Description=Calibre Offline Book Server
Documentation=https://manual.calibre-ebook.com/server.html
After=network-online.target srv-offline.mount
Wants=network-online.target srv-offline.mount
# Only start if the USB data drive is actually mounted
ConditionPathIsMountPoint=/srv/offline

[Service]
Type=simple
User=library
Group=library

# calibre-server does NOT auto-create a library — install.sh initializes it in step 3.
# The service will fail without an existing metadata.db; safe to start before books are synced.
ExecStart=/usr/bin/calibre-server \
    --port 8081 \
    --url-prefix /books \
    --listen-on 0.0.0.0 \
    /srv/offline/books/calibre-library

Restart=on-failure
RestartSec=10

MemoryMax=512M
TasksMax=64

[Install]
WantedBy=multi-user.target
