# systemd/srv-offline.mount
[Unit]
Description=SURVIVE offline data drive (/srv/offline)
# Short device wait so boot doesn't hang if USB is absent
DefaultDependencies=no
After=local-fs-pre.target
Before=local-fs.target

[Mount]
What=/dev/disk/by-label/survive-data
Where=/srv/offline
Type=ext4
Options=defaults,noatime,x-systemd.device-timeout=5s

[Install]
WantedBy=local-fs.target
