# systemd/survive-books.timer
[Unit]
Description=30-Minute Survive Books Sync Timer
Documentation=file:///srv/offline/scripts/sync/sync-books.sh

[Timer]
# Run every 30 minutes
OnCalendar=*:0/30

# If the Pi was offline at trigger time, run on next boot
Persistent=true

Unit=survive-books.service

[Install]
WantedBy=timers.target
