# systemd/survive-sync.service
[Unit]
Description=Survive Library Sync Agent
Documentation=file:///srv/offline/scripts/sync/sync-all.sh
# Requires internet — wait for a routable network connection
After=network-online.target srv-offline.mount
Wants=network-online.target srv-offline.mount
# sync-all.sh checks for connectivity and a mounted drive before downloading.

[Service]
Type=oneshot
# Run as library user; sudo is used only for service restarts via systemctl
User=library
Group=library

ExecStart=/srv/offline/scripts/sync/sync-all.sh

# Allow the sync to take a very long time (Wikipedia ZIM is ~100 GB)
TimeoutStartSec=0
TimeoutStopSec=600

# Log to journal; view with: journalctl -u survive-sync -f
StandardOutput=journal
StandardError=journal

# Prevent OOM kill during large downloads
OOMScoreAdjust=-100

[Install]
WantedBy=multi-user.target
